# deploy/docker/supervisord.conf
[supervisord]
nodaemon = true

[program:nginx]
command=/usr/sbin/nginx
user=root
autostart=true
autorestart = true
priority = 50

[program:php]
command=/usr/sbin/php-fpm8.2
user=root
autostart=true
autorestart = true
priority = 60
